# jndc_server/jndc-server-backend/src/main/resources/bin/jndc-server.service
[Unit]
Description=jndc-server
After=network.target
StartLimitIntervalSec=0

[Service]
User=app
Group=app
Type=forking
Restart=always
RestartSec=1
Environment=JAVA_HOME=/data/jdk8
WorkingDirectory=/data/app/jndc-server/bin
ExecStart=/data/app/jndc-server/bin/startup.sh
ExecReload=/data/app/jndc-server/bin/shutdown.sh
ExecStop=/data/app/jndc-server/bin/shutdown.sh
PrivateTmp=true
SuccessExitStatus=143

[Install]

WantedBy=multi-user.target